# portfolio-backend.service
[Unit]
Description=Rust Actix Backend Service
After=network.target

[Service]
WorkingDirectory=/srv/portfolio-backend
ExecStart=/srv/portfolio-backend/backend
Restart=on-failure
RestartSec=5s

Environment="DATABASE_URL=sqlite://portfolio.sqlite"

[Install]
WantedBy=multi-user.target
